# systemd/disklavier-web.service
[Unit]
Description=Disklavier Pi Web Interface
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/david/disklavier_pi
ExecStart=/home/david/disklavier_pi/venv/bin/uvicorn disklavier.api.app:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=5
User=david
Group=david
Environment=PATH=/home/david/disklavier_pi/venv/bin:/usr/local/bin:/usr/bin
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
