# systemd/fluidsynth-broadcast.service
[Unit]
Description=FluidSynth Piano Broadcast (DKC-800 → AirPlay)
After=pipewire-pulse.service
Wants=pipewire-pulse.service

[Service]
Type=simple
ExecStart=/usr/bin/fluidsynth -s -a pulseaudio -m alsa_seq -o midi.autoconnect=1 -g 1.0 /home/david/soundfonts/general_montage.sf2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
