# systemd/fluidsynth-broadcast.service
[Unit]
Description=FluidSynth Piano Broadcast (DKC-800 → AirPlay)
After=pipewire-pulse.service
Wants=pipewire-pulse.service

[Service]
Type=simple
ExecStart=/usr/bin/fluidsynth --server --no-shell --audio-driver=pulseaudio --midi-driver=alsa_seq -o midi.autoconnect=1 --gain=1.0 "/home/david/soundfonts/General Montage.sf2"
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
