# systemd/midi-autoconnect.service
[Unit]
Description=Auto-connect Network MIDI to DKC-800

[Service]
Type=oneshot
ExecStart=/home/david/disklavier_pi/scripts/midi-autoconnect.sh
